# claude-email.service
[Unit]
Description=Claude Email Agent
Documentation=https://github.com/cocodedk/claude-email
After=network-online.target claude-chat.service
Wants=network-online.target claude-chat.service

[Service]
Type=simple
WorkingDirectory=/home/cocodedk/0-projects/claude-email
EnvironmentFile=/home/cocodedk/0-projects/claude-email/.env
ExecStart=/home/cocodedk/0-projects/claude-email/.venv/bin/python3 /home/cocodedk/0-projects/claude-email/main.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=default.target
