# claude-email.service
[Unit]
Description=Claude Email Agent
Documentation=https://github.com/cocodedk/claude-email
After=network-online.target claude-chat.service
Wants=network-online.target claude-chat.service

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=__INSTALL_DIR__/.venv/bin/python3 __INSTALL_DIR__/main.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=default.target
